# ingestion_pipeline/cronschedule.crontab
# Run every hour
# 0 * * * * python3 /opt/application/src/main.py >> /var/log/flight_pipeline.log

# Every minute
* * * * * python3 /opt/application/src/main.py >> /var/log/flight_pipeline.log
